# packaging/patchbase-agent.timer
[Unit]
Description=Run PatchBase Agent Sync periodically

[Timer]
OnBootSec=2min
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
